# dcf.service
[Unit]
Description=Digital Crossover Filter

[Service]
ExecStart=/usr/bin/dcf

Restart=always
RestartSec=1

[Install]
WantedBy=default.target
